# etc/udev/rules.d/99-cd-audio-processing.rules
SUBSYSTEM=="block", KERNEL=="sr0", ACTION=="change", RUN+="/bin/systemctl start rip-audio-cd@sr0.service"
