# etc/systemd/system/rip-audio-cd@.service
[Unit]
Description=Rip audio CD on insertion

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/rip-audio-cd.sh %I
